# smokeping/rootfs/etc/services.d/fcgiwrap/run
#!/command/with-contenv bash
# Start fcgiwrap
exec fcgiwrap -f -s unix:/var/run/fcgiwrap.sock
