# smokeping/rootfs/etc/services.d/smokeping/run
#!/usr/bin/with-contenv bashio

exec /usr/sbin/smokeping --nodaemon